# zzzwifiHeart1
#!/bin/sh /etc/rc.common
sleep 190
killall heartbeat
/sbin/wifi down
