# zzzwifi2
#!/bin/sh /etc/rc.common
/sbin/wifi down
